# deploy/systemd/shuzhi-worker@.service
[Unit]
Description=星通数智 Task Worker %i
After=network-online.target postgresql.service redis.service
Wants=network-online.target

[Service]
Type=simple
User=shuzhi
Group=shuzhi
WorkingDirectory=/opt/shuzhi/backend
EnvironmentFile=/etc/shuzhi/shuzhi.env
Environment=WORKER_INSTANCE=%i
ExecStart=/opt/shuzhi/backend/.venv/bin/python -m scripts.task_worker
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
