# deploy/systemd/shuzhi-postgres-backup.service
[Unit]
Description=星通数智 PostgreSQL backup
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=shuzhi
Group=shuzhi
EnvironmentFile=/etc/shuzhi/shuzhi.env
ExecStart=/opt/shuzhi/deploy/scripts/shuzhi-postgres-backup.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
