# deploy/systemd/shuzhi-postgres-backup.timer
[Unit]
Description=Run 星通数智 PostgreSQL backup daily

[Timer]
OnCalendar=*-*-* 02:30:00
RandomizedDelaySec=1800
Persistent=true
Unit=shuzhi-postgres-backup.service

[Install]
WantedBy=timers.target
