# deploy/systemd/zhishu-worker@.service
[Unit]
Description=星通智数 Task Worker %i
After=network-online.target postgresql.service redis.service
Wants=network-online.target

[Service]
Type=simple
User=zhishu
Group=zhishu
WorkingDirectory=/opt/zhishu/backend
EnvironmentFile=/etc/zhishu/zhishu.env
Environment=WORKER_INSTANCE=%i
ExecStart=/opt/zhishu/backend/.venv/bin/python -m scripts.task_worker
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
