# deploy/systemd/zhishu-migrate.service
[Unit]
Description=星通智数 database migration
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=zhishu
Group=zhishu
WorkingDirectory=/opt/zhishu/backend
EnvironmentFile=/etc/zhishu/zhishu.env
ExecStart=/opt/zhishu/backend/.venv/bin/python -m scripts.db_migrate

[Install]
WantedBy=multi-user.target
